# linux_scripts/cron.d/check_quest_svn30
# Updates quests every 30 min
# MAILTO=""
*/30 * * * * root /home/samba/public/updates/update_quests_svn
