# linux_scripts/cron.d/check_cpu10
# Checks for overloaded processes
# MAILTO=""
*/10 * * * * root /home/samba/public/updates/check_load
