# linux_scripts/cron.d/start_server_daily
# Starts the EqEmu Server
MAILTO="angelox@nahunta.org"
03 04 * * * angelox /home/eqemu/server/run_server
