# trunk/linux_scripts/cron.d/char_bu_daily
# Backup accounts, characters, etc
# MAILTO=""
01 03 * * * root /home/samba/public/updates/backup_char
