# trunk/linux_scripts/cron.d/core_bu_daily
# Creates the new database
# MAILTO=""
31 03 * * * root /home/samba/public/updates/backup_main
